# packages/cloudron/supervisor/builder.conf
[program:builder]
priority=5
directory=/app/code
command=node /app/code/builder/apps/builder/server.js
user=cloudron
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
